# .config/systemd/user/RestartVSCodeServer.service
[Unit]
Description=Restart VSCode Server

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart VSCodeServer.service