# .config/systemd/user/UnknownHostname.service
[Unit]
Description=Fixes unknown hostname on boot


[Service]
ExecStart=/home/gloves/FixHostname.sh


[Install]
WantedBy=default.target
