# .config/systemd/user/VSCodeServer.service
[Unit]
Description=VSCode Server


[Service]
ExecStart=/home/gloves/VSCodeServer.sh
Type=simple
Restart=always


[Install]
WantedBy=default.target
